# kernel/usb_and_devices/g1/ubuntu_20.04_beta/saned/saned.service
[Unit]
Description=Saned ondemand

[Service]
#Type=oneshot
ExecStart=/opt/saned/config.sh

[Install]
WantedBy=multi-user.target

#/etc/systemd/system/saned.service
